# scripts/launchd/com.agentic-alert.daily.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agentic-alert.daily</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>__REPO_PATH__/scripts/run_daily_local.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <key>KeepAlive</key>
    <dict>
      <key>NetworkState</key>
      <true/>
    </dict>
    <key>WorkingDirectory</key>
    <string>__REPO_PATH__</string>
    <key>StandardOutPath</key>
    <string>__REPO_PATH__/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_PATH__/logs/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>UV_CACHE_DIR</key>
      <string>__REPO_PATH__/.cache/uv</string>
    </dict>
  </dict>
</plist>
